# crontab
30,45,59 12 * * * poetry run python jobs/check_weather.py  @@check_weather@@