# crontab
55 9 * * * poetry run python --version   @@demo@@