# crontab
35 11 * * * poetry run python jobs/check_weather.py  @@check_weather@@